# system/agent-c.service
[Unit]
Description=Jesse MCP Agent C - Risk Management
After=network.target

[Service]
Type=simple
User=bk
WorkingDirectory=/home/bk/source/jesse-mcp
ExecStart=/usr/bin/opencode run "Execute Agent C section from PRD_v1.0.0.md. Focus on risk analysis, portfolio optimization, and safety validation. Use monte_carlo_*, risk_report, risk_analyze_portfolio." -m zai-coding-plan/glm-4.5-flash
StandardOutput=append:/home/bk/source/jesse-mcp/logs/agent-c.log
StandardError=append:/home/bk/source/jesse-mcp/logs/agent-c.log
Restart=on-failure
RestartSec=30

[Install]
WantedBy=multi-user.target